# docker/spring/etc/supervisor/conf.d/envoy.conf
[program:envoy]
command=/usr/local/bin/envoy -c /etc/envoy.yaml --log-level info --log-path /var/log/envoy/access.log
autostart=true
autorestart=true
user=root
